# pkg/base/kovalromank-base.install
post_install() {
    post_upgrade
}

post_upgrade() {
    /usr/bin/systemctl enable polkit > /dev/null
    /usr/bin/systemctl enable NetworkManager > /dev/null
    /usr/bin/systemctl enable libvirtd > /dev/null
    /usr/bin/systemctl --user enable xdg-user-dirs-update > /dev/null
}

