# pkg/base/files/etc/profile.d/kovalromank-base
#!/usr/bin/env bash

export EDITOR='nvim'
export VISUAL='$EDITOR'

export PATH=$HOME/.local/bin:$PATH

alias ls=/usr/bin/exa

